# openresty/resource/docker-entrypoint.d/cron.d/crontab
01 00 * * * /bin/bash /docker-entrypoint.d/cron.d/nginx-logs.sh
